# trading_ai_monitor.service
[Unit]
Description=Trading AI Monitor (Bybit + Telegram)
After=network.target

[Service]
User=root
WorkingDirectory=/root/trading_ai_monitor_v2
ExecStart=/root/trading_ai_monitor_v2/venv/bin/python /root/trading_ai_monitor_v2/main.py

Restart=always
RestartSec=5

# Logs del bot irán a journalctl
StandardOutput=journal
StandardError=journal

# Ambiente
Environment="PYTHONUNBUFFERED=1"
Environment="TZ=America/Bogota"

[Install]
WantedBy=multi-user.target
