# trading_ai_monitor.service
[Unit]
Description=Trading AI Monitor - Monitoreo de Operaciones Abiertas
After=network.target

[Service]
WorkingDirectory=/home/usuario/trading_ai_monitor_v2
ExecStart=/home/usuario/trading_ai_monitor_v2/venv/bin/python3 /home/usuario/trading_ai_monitor_v2/main.py --mode monitor

Environment="PYTHONUNBUFFERED=1"

Restart=always
RestartSec=10

StandardOutput=append:/home/usuario/trading_ai_monitor_v2/logs/monitor.log
StandardError=append:/home/usuario/trading_ai_monitor_v2/logs/monitor_error.log

User=usuario
Group=usuario

[Install]
WantedBy=multi-user.target
